# deploy/elegir-cuidarme.service
[Unit]
Description=Elegir Cuidarme — Flask App
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/opt/elegir-cuidarme
ExecStart=/opt/elegir-cuidarme/venv/bin/gunicorn \
    --workers 2 \
    --bind 127.0.0.1:8000 \
    --access-logfile - \
    --error-logfile - \
    app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
